# scripts/templates/supervisord.conf
# Template supervisord.conf file.
# To generate, use `fab generate_supervisord_conf`
# Template variables:
#   project_dir, app_name, virtualenv_bin, user_name, processes (a dict of process counts)
# Defaults will be generated by the fab script, but can be overridden.

[supervisord]
childlogdir = {{ project_dir }}/logs
logfile = {{ project_dir }}/logs/supervisord.log
logfile_maxbytes = 5000000
logfile_backups = 4
loglevel = info
pidfile = /var/run/{{ app_name }}/supervisord.pid
umask = 022
environment=PYTHONUNBUFFERED="1",PATH="{{ path }}"

[unix_http_server]
file = /var/run/{{ app_name }}/supervisord.sock
chmod = 0700

[supervisorctl]
serverurl = unix:///var/run/{{ app_name }}/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:web]
command={{ virtualenv_bin }}/gunicorn twitter_feels.wsgi:application -c gunicorn.conf.py --access-logfile -
numprocs={{ processes.web }}
autostart=true
autorestart=true
stdout_logfile={{ project_dir }}/logs/web.log
stderr_logfile={{ project_dir }}/logs/web.error.log
user={{ user_name }}
directory={{ project_dir }}

[program:worker]
numprocs={{ processes.worker }}
process_name=%(program_name)s_%(process_num)02d
command={{ virtualenv_bin }}/python manage.py rqworker
autostart=true
autorestart=true
stdout_logfile={{ project_dir }}/logs/worker.log
stderr_logfile={{ project_dir }}/logs/worker.error.log
user={{ user_name }}
directory={{ project_dir }}

[program:stream]
command={{ virtualenv_bin }}/python manage.py stream
numprocs={{ processes.stream }}
autostart=true
autorestart=true
stdout_logfile={{ project_dir }}/logs/stream.log
stderr_logfile={{ project_dir }}/logs/stream.error.log
user={{ user_name }}
directory={{ project_dir }}

[program:scheduler]
command={{ virtualenv_bin }}/python manage.py rqscheduler
numprocs={{ processes.scheduler }}
autostart=true
autorestart=true
stdout_logfile={{ project_dir }}/logs/scheduler.log
stderr_logfile={{ project_dir }}/logs/scheduler.error.log
user={{ user_name }}
directory={{ project_dir }}

